# systemd/readsb-adsb.service.d/10-quiet-journal.conf
# Drop readsb's verbose per-message stdout out of the journal.
# Reversible: delete this file + `systemctl daemon-reload` + restart service.
# Does not affect readsb's network data feed used by the loggers.
[Service]
StandardOutput=null
StandardError=journal
